# images/kafka/resources/supervisor/kafka.conf
[program:kafka]
command=python /opt/kafka_2.9.2-0.8.2.1/pod/pod.py
